# autostart/waypaper.desktop
[Desktop Entry]
Type=Application
Name=waypaper
Exec=/usr/bin/waypaper --restore
Icon=system-run
StartupNotify=true
Terminal=false
